# config/puma.service
# ============================================================================
# SYSTEMD SERVICE FILE FOR PUMA WEB SERVER - MYCOWRITER
# ============================================================================
# This file defines the systemd service that manages Puma in production.
#
# INSTALLATION:
# sudo cp /opt/mycowriter/current/config/puma.service /etc/systemd/system/puma-mycowriter.service
# sudo systemctl daemon-reload
# sudo systemctl enable puma-mycowriter.service
#
# MANAGEMENT COMMANDS:
# - Start:   sudo systemctl start puma-mycowriter.service
# - Stop:    sudo systemctl stop puma-mycowriter.service
# - Restart: sudo systemctl restart puma-mycowriter.service
# - Status:  sudo systemctl status puma-mycowriter.service
# - Logs:    sudo journalctl -xeu puma-mycowriter.service
#
# AUTOMATIC DEPLOYMENT:
# Capistrano automatically restarts Puma after deployment via:
# config/deploy.rb: after "deploy:published", "systemd_puma:restart"
#
# ⚠️  CRITICAL: Changes to this file require:
# 1. Copy to /etc/systemd/system/puma-mycowriter.service
# 2. sudo systemctl daemon-reload
# 3. sudo systemctl restart puma-mycowriter.service
# ============================================================================

[Unit]
Description=Puma HTTP Server for MycoWriter (production)
After=network.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/opt/mycowriter/current

# ============================================================================
# ENVIRONMENT VARIABLES - CRITICAL FOR RAILS OPERATION
# ============================================================================
Environment=RAILS_ENV=production
Environment=RBENV_ROOT=/home/deploy/.rbenv
Environment=PATH=/home/deploy/.rbenv/shims:/home/deploy/.rbenv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# ============================================================================
# PROCESS MANAGEMENT
# ============================================================================
ExecStartPre=/bin/rm -f /opt/mycowriter/shared/tmp/sockets/puma.sock
ExecStartPre=/bin/rm -f /opt/mycowriter/shared/tmp/pids/puma.pid
ExecStartPre=/bin/rm -f /opt/mycowriter/shared/tmp/pids/puma.state
ExecStart=/home/deploy/.rbenv/shims/bundle exec puma -C /opt/mycowriter/current/config/puma.rb
ExecReload=/bin/kill -USR1 $MAINPID
ExecStop=/bin/kill -SIGTERM $MAINPID

# ============================================================================
# RESTART POLICY - ENSURES HIGH AVAILABILITY
# ============================================================================
Restart=always
RestartSec=10
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
StartLimitIntervalSec=300
StartLimitBurst=5

# ============================================================================
# LOGGING
# ============================================================================
StandardOutput=append:/opt/mycowriter/shared/log/puma_stdout.log
StandardError=append:/opt/mycowriter/shared/log/puma_stderr.log
SyslogIdentifier=puma-mycowriter

[Install]
WantedBy=multi-user.target
